# systemd/cod-confirm.service
[Unit]
Description=Glitch COD Confirm \u2014 voice AI COD confirmation service
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=support
WorkingDirectory=/home/support/glitch-cod-confirm
EnvironmentFile=/home/support/glitch-cod-confirm/.env
Environment=PATH=/home/support/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/support/.npm-global/bin/pnpm start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cod-confirm

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/support/glitch-cod-confirm

[Install]
WantedBy=multi-user.target
